# .config/systemd/user/dunst.service
[Unit]
Description=Dunst
[Service]
Type=simple
Restart=always
RestartSec=30
TimeoutStartSec=0
ExecStart=dunst
[Install]
WantedBy=basic.target
