# .config/systemd/user/inhousenet.service
[Unit]
Description=Inhouse Network Monitor
[Service]
Type=oneshot
ExecStart=/home/victor/.local/bin/scripts/inhousenet
